# ansible/template/app.service.j2
# ansible/templates/app.service.j2
[Unit]
Description=Backend Application
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/app
ExecStart=/usr/local/bin/docker-compose -f /opt/app/docker-compose.uat.yml up -d
ExecStop=/usr/local/bin/docker-compose -f /opt/app/docker-compose.uat.yml down
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target
